# linux/.config/autostart/fluxgui.desktop
[Desktop Entry]
Type=Application
Name=f.lux indicator applet
Exec=sh -c "\\"\\$SHELL\\" -c fluxgui"
Icon=fluxgui
X-GNOME-Autostart-enabled=true

